# etc/systemd/system/unalix.socket
[Unit]
Description=Unalix socket

[Socket]
ListenStream=/run/unalix.sock

[Install]
WantedBy=sockets.target
